# templates/borgmatic-backup.timer.j2
[Unit]
Description=Run borgmatic-backup.service {{ calendar }}

[Timer]
OnCalendar={{ calendar }}
