# templates/borgmatic-backup.service.j2
# {{ ansible_managed }}
# DO NOT EDIT THIS FILE

[Unit]
Description=Borgmatic backup

[Service]
Type=oneshot
Environment="BORG_RSH={{ borg_rsh }}"
Environment="BORG_KEYS_DIR=/etc/borgmatic"
{% if run_before_job %}
ExecStartPre={{ run_before_job }}
{% endif %}
ExecStart={{ borgmatic_binary }} -v 2
