# templates/borgmatic-backup.timer.j2
# {{ ansible_managed }}
# DO NOT EDIT THIS FILE

[Unit]
Description=Run borgmatic-backup.service {{ calendar }}

[Timer]
OnCalendar={{ calendar }}

[Install]
WantedBy=timers.target
